# config/systemd/user/mailagent.service
[Unit]
Description=Sync mail
RefuseManualStart=no
RefuseManualStop=yes

[Service]
Type=oneshot
ExecStart=/usr/local/bin/mailagent -aV
